# shell/crontab/crontab秒级.sh
#!/bin/bash
# --------------------------------------------------
#Author:    LJ
#Email:		admin@attacker.club


#Last Modified: 2018-02-24 18:03:41
#Description:	crontab秒级计划任务
# --------------------------------------------------


* * * * * sleep 10; php /home/php/crontab/tolog.php