# shell/crontab/while秒级任务.sh
#!/bin/bash
# --------------------------------------------------
#Author:	LGhost
#Email:		admin@attacker.club
#Site:		blog.attacker.club

#Last Modified: 2017-12-15 23:17:15
#Description:	 死循环秒级任务
# --------------------------------------------------

while true do 
command 
sleep 20 ###睡眠20秒执行一次 
done

#while true ;do
#/home/aide.sh >/dev/null 2>&1
#sleep 7200 #2小时检测
#done