# shell/crontab/while秒级任务.sh
#!/bin/bash
# --------------------------------------------------
#Author:    LJ
#Email:		admin@attacker.club


#Last Modified: 2018-02-24 18:03:40
#Description:	 死循环秒级任务
# --------------------------------------------------

while true do 
command 
sleep 20 ###睡眠20秒执行一次 
done

#while true ;do
#/home/aide.sh >/dev/null 2>&1
#sleep 7200 #2小时检测
#done